# etc/systemd/system/msq23.service
[Unit]
Description=Matheus S Queiroga Service
After=network.services

[Service]
ExecStart=/bin/Sh23Start.sh
Restart=always

